# deploy/gcp/systemd/tradegy-mes-entry.timer
[Unit]
Description=Run Tradegy MES 0DTE paper entry on market sessions

[Timer]
OnCalendar=Mon..Thu *-*-* 10:00:00 America/New_York
Persistent=false
RandomizedDelaySec=0
AccuracySec=1s
Unit=tradegy-mes-entry.service

[Install]
WantedBy=timers.target
